# deploy/qc-pii-gate.service.d/v11int8.conf
# Drop-in: repoint ossredact-gate to the v11r5-base ONNX-INT8 CPU tier (gate_service_cpu.py).
# Replaces the v7 OpenVINO/NPU OVTier. Revert = remove this file + daemon-reload + restart
# (the base unit's original ExecStart -> gate_service.py / OVTier is preserved).
[Service]
ExecStart=
ExecStart=/opt/ossredact-npu/.venv/bin/python /opt/ossredact-npu/gate_service_cpu.py
Environment=CPU_GATE_PORT=8001
Environment=CPU_GATE_MODEL=/opt/ossredact-npu/model-v11r5-base-int8
Environment=CUDA_VISIBLE_DEVICES=
