# deploy/qc-pii-gate-cpu.service
[Unit]
Description=ossredact CPU PII gate (xlm-r-base v11r5 int8 on CPU) - own-use redaction sidecar :8011
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ossredact
WorkingDirectory=/opt/ossredact-gpu-gate
# CPU-only: explicitly hide all GPUs so this never contends with the GPU gate on the GPU
Environment=CUDA_VISIBLE_DEVICES=
Environment=CPU_GATE_PORT=8011
Environment=CPU_GATE_MODEL=/opt/ossredact/models/privacy-filters/pii-gpu-xlmr-base-v11r5
ExecStart=/opt/ossredact/.venv-pii/bin/python /opt/ossredact-gpu-gate/gate_service_cpu.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
