# deploy/qc-pii-gate-cpu.service
[Unit]
Description=qc-pii CPU PII gate (xlm-r-base v11r5 int8 on CPU) - own-use redaction sidecar :8011
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=steven
WorkingDirectory=/home/steven/sparx-gpu-gate
# CPU-only: explicitly hide all GPUs so this never contends with the GPU gate on card 4
Environment=CUDA_VISIBLE_DEVICES=
Environment=CPU_GATE_PORT=8011
Environment=CPU_GATE_MODEL=/home/steven/Sparx/models/privacy-filters/pii-gpu-xlmr-base-v11r5
ExecStart=/home/steven/Sparx/.venv-pii/bin/python /home/steven/sparx-gpu-gate/gate_service_cpu.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
