# deploy/ossredact-gate.service.d/v11int8.conf
# Drop-in: repoint the gate to the v11r5-base ONNX-INT8 CPU tier (gate_service_cpu.py).
# Replaces a prior OpenVINO/NPU tier. Revert = remove this file + daemon-reload + restart
# (the base unit's original ExecStart is preserved).
[Service]
ExecStart=
ExecStart=/opt/ossredact/.venv/bin/python /opt/ossredact/gate_service_cpu.py
Environment=CPU_GATE_PORT=8001
Environment=CPU_GATE_MODEL=/opt/ossredact/models/pii-xlmr-base-int8
Environment=CUDA_VISIBLE_DEVICES=
